# deployment/analytics-engine.service
[Unit]
Description=Analytics Engine FastAPI Service
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/analytics-engine
Environment="PATH=/var/www/analytics-engine/venv/bin"
ExecStart=/var/www/analytics-engine/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8001

[Install]
WantedBy=multi-user.target
