# deployment/cdn-worker.conf
[program:cdn-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/metapilot/artisan queue:work --queue=cdn-ingestion --sleep=1 --tries=3 --timeout=60
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/metapilot/storage/logs/cdn-worker.log
stopwaitsecs=3600
