# deployment/zombie-cleanup-worker.conf
[program:metapilot-cleanup-worker]
process_name=%(program_name)s
command=/bin/bash -c "while true; do /var/www/metapilot/deployment/zombie_cleanup.sh; sleep 900; done"
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/metapilot/storage/logs/cleanup-worker.log
stopwaitsecs=3600
